# roles/virtualization/files/google-cloud-sdk/rootfs/etc/profile.d/google-cloud-sdk.sh
if test -n "${BASH_VERSION}"; then
    . /opt/google/google-cloud-sdk/path.bash.inc
    . /opt/google/google-cloud-sdk/completion.bash.inc
fi
